# files/rsyncd-clean.cron
SHELL=/bin/bash
PATH=/sbin:/bin:/usr/sbin:/usr/bin
MAILTO=root
55 2 * * * root /usr/bin/rsyncd-clean.sh

